# build_dir/target-mips-rlx5281-linux/root-realtek/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmpproxy
/usr/sbin/igmpproxy
/etc/config/igmpproxy
