# build_dir/target-mips-rlx5281-linux/root-realtek/usr/lib/opkg/info/wshaper.list
/etc/init.d/wshaper
/usr/sbin/wshaper.htb
/etc/config/wshaper
